# config/systemd/user/time-wallpaper.service
[Unit]
Description=Update wallpaper based on current time

[Service]
Type=simple
ExecStart=/home/yog/bin/time-wallpaper
